# rbp-timelapse/rc.local
#!/bin/bash

sudo -u pi /opt/raspberry-pi/autostart.sh

exit 0

